# deploy/systemd/kinder-vision-api.service
[Unit]
Description=Kinder Vision API (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/kinder-vision
EnvironmentFile=-/etc/kinder-vision/kinder-vision.env
# Fallback defaults if EnvironmentFile not present:
Environment=KINDER_MEMORY_DIR=/var/lib/kinder-vision/memory
Environment=KINDER_TMP_DIR=/var/lib/kinder-vision/tmp
Environment=KINDER_REPORTS_DIR=/var/lib/kinder-vision/reports
Environment=KINDER_TASK_TTL_SEC=86400
# Optional AI settings:
# Environment=KINDER_AI_API_KEY=...
# Environment=KINDER_AI_BASE_URL=https://api.openai.com/v1
# Environment=KINDER_AI_MODEL=gpt-4o-mini

ExecStart=/opt/kinder-vision/.venv/bin/uvicorn src.api:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
