# deploy/systemd/kinder-vision-worker.service
[Unit]
Description=Kinder Vision Worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/kinder-vision
EnvironmentFile=-/etc/kinder-vision/kinder-vision.env
# Fallback defaults if EnvironmentFile not present:
Environment=KINDER_MEMORY_DIR=/var/lib/kinder-vision/memory
Environment=KINDER_TMP_DIR=/var/lib/kinder-vision/tmp
Environment=KINDER_REPORTS_DIR=/var/lib/kinder-vision/reports
# Optional AI settings:
# Environment=KINDER_AI_API_KEY=...
# Environment=KINDER_AI_BASE_URL=https://api.openai.com/v1
# Environment=KINDER_AI_MODEL=gpt-4o-mini

ExecStart=/home/ubuntu/miniconda/envs/kinder-vision-py311/bin/python -m src /data/media/input.mp4 --stride 4 --pose pose
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
